# src/autostart/bbrun.desktop
[Desktop Entry]
Encoding=UTF-8
Type=Application
Name=BlackBox Run Command
GenericName=bbrun 1.6
Comment=Run Command DockApp
Icon=system-run
OnlyShowIn=BLACKBOX;
TryExec=/usr/bin/bbrun
Exec=bbrun -a
Terminal=false
Categories=DockApp;WMAKER;
StartupNotify=false
StartupWMClass=bbrun
X-GNOME-Autostart-enabled=true
