# src/autostart/compton.desktop
[Desktop Entry]
Version=1.0
Type=Application
Name=Compton
GenericName=X compositor
Comment=An X compositor
Categories=Utility;
TryExec=compton
Exec=compton --config ~/.config/compton/compton.conf
Icon=xcompmgr
OnlyShowIn=FLUXBOX;BLACKBOX;OPENBOX;ICEWM;WMAKER;FVWM;JWM;PEKWM;ADWM;WAIMEA;
Hidden=true
# Thanks to quequotion for providing this file!
X-Arch-Package=compton-git 0.1_beta2.59.g23d1dd1-1
