# src/autostart/docker.desktop
[Desktop Entry]
Encoding=UTF-8
Type=Application
Name=System Tray Docker
GenericName=System Tray
Comment=Provide a system tray in the dock
Icon=GNUstep3D
#only necessary in WMs that do not have a systray
OnlyShowIn=BLACKBOX;WMAKER;PEKWM;
Exec=docker -wmaker -iconsize 16
TryExec=docker
Terminal=false
Categories=DockApp;WMAKER;
StartupNotify=false
#use wmsystemtray instead
Hidden=true
