# src/autostart/retrovol.desktop
[Desktop Entry]
Type=Application
Name=Retrovol sound mixer
GenericName=Retrovol
Comment=Retrovol sound mixer
Icon=wmaker-sound_mixer.png
TryExec=/usr/bin/retrovol
Exec=retrovol -hide
Hidden=true
NotShowIn=GNOME;KDE;XFCE;LXDE;
Categories=TrayIcon;GTK;
StartupNotify=false
X-AppInstall-Package=retrovol
X-AppInstall-Section=custom
X-Arch-Package=retrovol 0.14.1-1
X-Arch-Package-Repository=custom
X-Arch-Package-Name=retrovol
X-Arch-Package-Version=0.14.1-1
X-Arch-Package-Description=A sound volume mixer
X-Arch-Package-Groups=None
X-Arch-Package-URL=https://github.com/pizzasgood/retrovol
X-GNOME-Autostart-enable=true
