# src/autostart/wmbattery.desktop
[Desktop Entry]
Encoding=UTF-8
Type=Application
Name=Battery Dock App
GenericName=Monitor
Comment=A batter monitor dock app
Icon=wmaker-dock
OnlyShowIn=FLUXBOX;OPENBOX;LXDE;BLACKBOX;
Exec=wmbattery
TryExec=/usr/bin/wmbattery
Terminal=false
Categories=DockApp;WMAKER;
StartupNotify=false
