# src/autostart/wmCalClock.desktop
[Desktop Entry]
Encoding=UTF-8
Type=Application
Name=Calendar Clock Dock App
GenericName=Clock
Comment=Calendar clock dock app
Icon=GNUstep3D
OnlyShowIn=FLUXBOX;OPENBOX;LXDE;BLACKBOX;PEKWM;WMAKER;
Exec=wmCalClock
TryExec=wmCalClock
Terminal=false
Categories=DockApp;WMAKER;
StartupNotify=false
