# src/autostart/wmcpuload.desktop
[Desktop Entry]
Encoding=UTF-8
Type=Application
Name=CPU Load Dock App
GenericName=Load Monitor
Comment=Show the CPU load
Icon=GNUstep3D
OnlyShowIn=FLUXBOX;OPENBOX;LXDE;BLACKBOX;PEKWM;WMAKER;
Exec=wmcpuload
TryExec=wmcpuload
Terminal=false
Categories=DockApp;WMAKER;
StartupNotify=false
