# src/autostart/wmnet.desktop
[Desktop Entry]
Encoding=UTF-8
Type=Application
Name=Network Monitor
GenericName=Monitor
Comment=Monitor network interfaces in a dock app
Icon=wmaker-dock
OnlyShowIn=FLUXBOX;OPENBOX;LXDE;BLACKBOX;
Exec=wmnet -w -W eth0 -x 1000000000 -l
TryExec=/usr/bin/wmnet
Terminal=false
Categories=DockApp;WMAKER;
StartupNotify=false
