# src/autostart/wmsysmon.desktop
[Desktop Entry]
Encoding=UTF-8
Type=Application
Name=System Monitor
GenericName=Monitor
Comment=Monitor system in a dock app
Icon=GNUstep3D
OnlyShowIn=FLUXBOX;OPENBOX;LXDE;BLACKBOX;PEKWM;ADWM;
Exec=wmsysmon
TryExec=/usr/bin/wmsysmon
Terminal=false
Categories=DockApp;WMAKER;
StartupNotify=false
