# src/autostart/wmtimer.desktop
[Desktop Entry]
Encoding=UTF-8
Type=Application
Name=Timer Dock App
GenericName=Timer
Comment=Timer dock app
Icon=wmaker-dock
OnlyShowIn=FLUXBOX;OPENBOX;LXDE;BLACKBOX;
Exec=wmtimer
TryExec=/usr/bin/wmtimer
Terminal=false
Categories=DockApp;WMAKER;
StartupNotify=false
