# src/autostart/wmweather+.desktop
[Desktop Entry]
Encoding=UTF-8
Type=Application
Name=Weather Monitor Dock App Plus
GenericName=Weather
Comment=Monitor weather conditions
Icon=wmaker-dock
OnlyShowIn=FLUXBOX;OPENBOX;LXDE;BLACKBOX;
Exec=wmweather+ -s CYEG
TryExec=/usr/bin/wmweather+
Terminal=false
Categories=DockApp;WMAKER;
StartupNotify=false
