# src/autostart/xde-panel.desktop
[Desktop Entry]
Type=Application
Name=XDE Panel
GenericName=Panel
Comment=Panel for X Desktop Environment.
Icon=perlpanel-panel
TryExec=/usr/bin/xde-panel
Exec=sh -c "xde-panel $HOME/.config/xde-panel/${XDG_CURRENT_DESKTOP:-default}/rc.xml"
OnlyShowIn=FLUXBOX;BLACKBOX;OPENBOX;WMAKER;PEKWM;ADWM;WAIMEA;ECHINUS;SCROTWM;
Categories=System;Utility;XDE;
StartupNotify=true
StartupWMClass=xde-panel
X-AppInstall-Package=xde-panel
